# netlink_notify/files/netlink_notify_fix.init
#!/bin/sh /etc/rc.common
START=99

start()
{
	/usr/sbin/netlink_notify_check &
}

stop()
{
	kill -9  $(ps | grep -v grep | grep netlink_notify_check | awk '{print $1}')
}

restart()
{
	stop;
	start;
}
