# netlink_notify/files/netlink_notify.init
#!/bin/sh /etc/rc.common
START=99

start()
{
	/usr/sbin/netlink_notify_app &
}

stop()
{
	killall netlink_notify_app
}

restart()
{
	stop;
	start;
}
