# local-luci/files/root/etc/uci-defaults/build_version
#!/bin/sh

BUILD_DATE=$(sed -n '9p' /rom/etc/banner)
echo $BUILD_DATE > /etc/build_version
