# local-luci/files/root/etc/uci-defaults/hw_type
#!/bin/sh

HW_TYPE=$(cat /etc/hw_type | awk -F '=' '{print $1}' | awk -F '_' '{print $6" "$7}')

echo $HW_TYPE > /etc/hw_type
